# scripts/launch-agents/run.factory.f1-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>run.factory.f1-daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/gaganarora/Desktop/my projects/active_graph/scripts/f1-daemon.mjs</string>
    <string>--interval-seconds</string>
    <string>300</string>
    <string>--auto-respawn</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/gaganarora/Desktop/my projects/active_graph</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/gaganarora/.factory/f1-daemon.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gaganarora/.factory/f1-daemon.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
