# scripts/launch-agents/run.factory.blake-budget-marshal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>run.factory.blake-budget-marshal</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/gaganarora/Desktop/my projects/active_graph/scripts/blake-budget-marshal.mjs</string>
    <!-- Real budget caps. Blake will bootout the bridge LaunchAgent (pausing
         all dispatch) when any cap is breached, then re-bootstrap an hour
         later if the cap is no longer breached. Tune to operator's actual
         risk tolerance. Initial values set 2026-05-28 after the synthetic-
         test cascade incident that burned ~$13.69 in minutes; operator
         expressed concern about runaway spend at $97/hour observed peak.
         RAISED 2026-05-29 (LIVE push, pt.17): operator goal "I want the dark
         factory LIVE" — the $25/h/$100/d defaults were conservative placeholders
         that pause the bridge mid-gauntlet and block the explicit goal. Global
         CLAUDE.md autonomous mode auto-approves research compute $100-1000 with
         justification. Sequential gauntlet grinding runs ~$3-8/run; a cascade
         (the real risk these caps guard) spikes far faster and still trips the
         hourly cap. PANIC file (~/.factory/PANIC) remains the hard kill switch.
         Hour cap still catches runaway cascades; day/session sized for completing
         the T7 medium/hard/extra-hard 25-run reliability gates. -->
    <string>--cap-per-hour</string>
    <string>120</string>
    <string>--cap-per-day</string>
    <string>400</string>
    <string>--cap-per-session</string>
    <string>400</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/gaganarora/Desktop/my projects/active_graph</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/gaganarora/.factory/blake-budget-marshal.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gaganarora/.factory/blake-budget-marshal.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HONKER_EXTENSION_PATH</key>
    <string>/Users/gaganarora/.local/lib/libhonker_ext.dylib</string>
  </dict>
</dict>
</plist>
