# scripts/launch-agents/run.factory.safety-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>run.factory.safety-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/gaganarora/Desktop/my projects/active_graph/scripts/safety-monitor.mjs</string>
  </array>
  <!-- Default: stub harm-check (offline regex, zero token cost — catches
       secrets, destructive shell, exfiltration, DROP TABLE, force-push, etc.).
       Add a "--live" <string> arg above to ALSO run model-based harm review on
       Claude Code auth (opus-4.8, no API key). --live costs ~one cheap claude
       call per flywheel.diff.proposed. -->


  <key>WorkingDirectory</key>
  <string>/Users/gaganarora/Desktop/my projects/active_graph</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/gaganarora/.factory/safety-monitor.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gaganarora/.factory/safety-monitor.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HONKER_EXTENSION_PATH</key>
    <string>/Users/gaganarora/.local/lib/libhonker_ext.dylib</string>
  </dict>
</dict>
</plist>
